# debian/prometheus-home-sizes-daily.timer
[Unit]
Description=Daily task updating sizes of subdirectories
Requires=prometheus-home-sizes.service

[Timer]
Unit=prometheus-home-sizes.service
OnCalendar=*-*-* 5:00:00
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 20:00:00

[Install]
WantedBy=multi-user.target